# roles/ironic_http/templates/ironic-http.service.j2
[Unit]
Description=Ironic HTTP Server Service (IPA Images and Boot ISOs)
Documentation=https://docs.openstack.org/ironic/latest/
After=docker.service ironic-network.service
Requires=docker.service ironic-network.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Run as unprivileged user
User={{ ironic_service_user }}
Group={{ ironic_service_group }}

# Pull the latest image
ExecStartPre=-{{ docker_bin }} pull {{ ironic_image_repo }}/ironic-httpd:{{ ironic_image_tag }}

# Ensure images directory exists
ExecStartPre=-mkdir -p {{ ironic_http_images_dir }}

# Start the Ironic HTTP container
ExecStartPre=-{{ docker_bin }} rm -f {{ ironic_http_container_name }}
ExecStart={{ docker_bin }} run --name {{ ironic_http_container_name }} \
    -d \
    --restart=always \
    --network ironic-network \
    -v {{ ironic_http_images_dir }}:/var/www/html:ro,z \
    -p {{ ironic_http_port }}:8080 \
    {{ ironic_image_repo }}/ironic-httpd:{{ ironic_image_tag }}

# Graceful shutdown
ExecStop={{ docker_bin }} stop {{ ironic_http_container_name }}
ExecStopPost={{ docker_bin }} rm {{ ironic_http_container_name }}

TimeoutStartSec=60
TimeoutStopSec=30

# Security
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
