# roles/mariadb/templates/mariadb.service.j2
[Unit]
Description=Ironic MariaDB Database Service
After=docker.service
Requires=docker.service

[Service]
Type=simple
Restart=on-failure
RestartSec=10s

ExecStartPre=-{{ docker_bin }} stop {{ mariadb_container_name }}
ExecStartPre=-{{ docker_bin }} rm {{ mariadb_container_name }}
ExecStartPre=-{{ docker_bin }} pull {{ mariadb_image }}

ExecStart={{ docker_bin }} run --name {{ mariadb_container_name }} \
    --rm \
    --network {{ docker_network_name }} \
    --env-file {{ ironic_etc_dir }}/mariadb.env \
    -v {{ mariadb_data_dir }}:/var/lib/mysql:z \
    -p {{ mariadb_bind_addr }}:{{ mariadb_port }}:3306 \
    {{ mariadb_image }}

ExecStop={{ docker_bin }} stop {{ mariadb_container_name }}

TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
